# hammertalk.service
[Unit]
Description=Push-to-talk transcription daemon for Wayland
After=pipewire.service

[Service]
ExecStart=%h/.local/bin/hammertalk
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
